# dot_config/systemd/user/dotfiles-improve.service
[Unit]
Description=dotfiles self-improvement run

[Service]
Type=oneshot
ExecStart=%h/bin/dotfiles-improve.sh
